# daemon/sovnas-daemon.service
[Unit]
Description=sovnas host daemon for Urbit NAS
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
# Replace YOUR-SHIP below with your ship name (e.g. sampel-palnet)
ExecStart=/usr/bin/python3 /opt/sovnas/sovnas-daemon.py --pier /media/data/docker/volumes/YOUR-SHIP/_data/YOUR-SHIP --root /home/nativeplanet/sovnas
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
