# daemon/sovnas-daemon.service
[Unit]
Description=sovnas host daemon for Urbit NAS
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
# Copy sovnas.config.template.json to /opt/sovnas/sovnas.config.json and edit it
ExecStart=/usr/bin/python3 /opt/sovnas/sovnas-daemon.py --config /opt/sovnas/sovnas.config.json
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
